# sandbox/supervisord-replit.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=Xvfb :1 -screen 0 1280x1029x24
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=10

[program:chrome]
command=chromium
    --display=:1
    --window-size=1280,1029
    --start-maximized
    --no-sandbox
    --disable-dev-shm-usage
    --disable-setuid-sandbox
    --disable-accelerated-2d-canvas
    --disable-gpu
    --no-first-run
    --no-default-browser-check
    --disable-infobars
    --disable-popup-blocking
    --disable-gpu-sandbox
    --no-xshm
    --disable-notifications
    --disable-extensions
    --disable-web-security
    --user-data-dir=/tmp/chrome-data
    --remote-debugging-address=0.0.0.0
    --remote-debugging-port=8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=20
startretries=3
startsecs=5

[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=2

[program:x11vnc]
command=x11vnc -display :1 -nopw -shared -listen 0.0.0.0 -xkb -forever -rfbport 5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=40
startsecs=3

[program:websockify]
command=websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=45
startsecs=3

[program:app]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/home/runner/workspace/sandbox
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50

[group:services]
programs=xvfb,chrome,socat,x11vnc,websockify,app
